# deploy/statusbot.service
[Unit]
Description=StatusBot (Kotlin)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=statusbot
Group=statusbot
WorkingDirectory=/opt/statusbot
EnvironmentFile=/etc/statusbot/statusbot.env
ExecStart=/usr/bin/java -Xms16m -Xmx80m -XX:+UseSerialGC -XX:MaxMetaspaceSize=64m -XX:ReservedCodeCacheSize=32m -XX:+ExitOnOutOfMemoryError -jar /opt/statusbot/statusbot.jar
SuccessExitStatus=143
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
